# docker/etc/services.d/10-son/run
#!/usr/bin/execlineb -P
/volume/bin/server