# docker/etc/services.d/10-son/run
#!/usr/bin/execlineb -P
exec /volume/go/bin/server